# deploy_aws_terraform_singleserver/ansible/templates/routerdb.socket.j2
[Unit]
Description=routerdb socket

[Socket]
ListenStream=/run/routerdb.sock

[Install]
WantedBy=sockets.target